# raspberrypi/etc/systemd/system/raspberrypi_exporter.timer
[Unit]
Description=Timer for the Raspberry Pi Prometheus exporter

[Timer]
OnCalendar=*:*:0,15,30,45

[Install]
WantedBy=timers.target
